# controller_yoga/files/fetch-crl.cron
45 */6 * * * root /usr/sbin/fetch-crl -q -r 360
